# qubes-core-qrexec/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/pam.d/
-rw-r--r--	root/root	etc/pam.d/qrexec
drwxr-xr-x	root/root	etc/qubes-rpc/
lrwxrwxrwx	root/root	etc/qubes-rpc/policy.Ask -> /var/run/qubes/policy-agent.sock
lrwxrwxrwx	root/root	etc/qubes-rpc/policy.Notify -> /var/run/qubes/policy-agent.sock
drwxr-xr-x	root/root	etc/qubes/
drwxr-xr-x	root/root	etc/qubes/rpc-config/
-rw-r--r--	root/root	etc/qubes/rpc-config/README
drwxr-xr-x	root/root	etc/xdg/
drwxr-xr-x	root/root	etc/xdg/autostart/
-rw-r--r--	root/root	etc/xdg/autostart/qrexec-policy-agent.desktop
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/qrexec-client-vm
-rwxr-xr-x	root/root	usr/bin/qrexec-fork-server
lrwxrwxrwx	root/root	usr/bin/qrexec-policy -> qrexec-policy-exec
-rwxr-xr-x	root/root	usr/bin/qrexec-policy-agent
-rwxr-xr-x	root/root	usr/bin/qrexec-policy-daemon
-rwxr-xr-x	root/root	usr/bin/qrexec-policy-exec
-rwxr-xr-x	root/root	usr/bin/qrexec-policy-graph
-rwxr-xr-x	root/root	usr/bin/qrexec-policy-restore
-rwxr-xr-x	root/root	usr/bin/qubes-policy
-rwxr-xr-x	root/root	usr/bin/qubes-policy-admin
drwxr-xr-x	root/root	usr/include/
-rw-r--r--	root/root	usr/include/libqrexec-utils.h
-rw-r--r--	root/root	usr/include/qrexec.h
drwxr-xr-x	root/root	usr/lib/
lrwxrwxrwx	root/root	usr/lib/libqrexec-utils.so -> libqrexec-utils.so.3
-rwxr-xr-x	root/root	usr/lib/libqrexec-utils.so.3
drwxr-xr-x	root/root	usr/lib/python3.9/
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec-4.1.15-py3.9.egg-info/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec-4.1.15-py3.9.egg-info/PKG-INFO
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec-4.1.15-py3.9.egg-info/SOURCES.txt
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec-4.1.15-py3.9.egg-info/dependency_links.txt
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec-4.1.15-py3.9.egg-info/top_level.txt
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__init__.py
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/__init__.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/__init__.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/client.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/client.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/exc.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/exc.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/server.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/server.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/utils.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/__pycache__/utils.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/client.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/exc.py
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/glade/
-rw-rw-r--	root/root	usr/lib/python3.9/site-packages/qrexec/glade/RPCConfirmationWindow.glade
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/policy/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__init__.py (EMPTY)
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/__init__.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/__init__.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/admin.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/admin.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/admin_client.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/admin_client.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/parser.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/parser.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/parser_compat.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/parser_compat.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/utils.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/__pycache__/utils.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/admin.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/admin_client.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/parser.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/parser_compat.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/policy/utils.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/server.py
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/tests/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__init__.py (EMPTY)
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/__init__.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/__init__.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/cli.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/cli.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/gtkhelpers.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/gtkhelpers.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/policy_admin.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/policy_admin.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/policy_cache.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/policy_cache.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/policy_parser.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/policy_parser.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/qrexec_policy_daemon.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/qrexec_policy_daemon.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/rpcconfirmation.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/rpcconfirmation.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/server.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/__pycache__/server.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/cli.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/gtkhelpers.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/policy_admin.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/policy_cache.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/policy_parser.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/qrexec_policy_daemon.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/rpcconfirmation.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/server.py
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__init__.py (EMPTY)
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/__init__.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/__init__.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/agent.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/agent.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/daemon.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/daemon.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/qrexec.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/qrexec.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/util.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/__pycache__/util.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/agent.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/daemon.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/qrexec.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tests/socket/util.py
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/tools/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__init__.py (EMPTY)
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/__init__.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/__init__.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qrexec_policy_agent.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qrexec_policy_agent.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qrexec_policy_daemon.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qrexec_policy_daemon.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qrexec_policy_exec.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qrexec_policy_exec.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qrexec_policy_graph.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qrexec_policy_graph.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qrexec_policy_restore.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qrexec_policy_restore.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qubes_policy.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qubes_policy.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qubes_policy_admin.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/__pycache__/qubes_policy_admin.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/qrexec_policy_agent.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/qrexec_policy_daemon.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/qrexec_policy_exec.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/qrexec_policy_graph.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/qrexec_policy_restore.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/qubes_policy.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/tools/qubes_policy_admin.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qrexec/utils.py
drwxr-xr-x	root/root	usr/lib/qubes/
-rwxr-xr-x	root/root	usr/lib/qubes/qrexec-agent
lrwxrwxrwx	root/root	usr/lib/qubes/qrexec-client-vm -> ../../bin/qrexec-client-vm
-rwxr-xr-x	root/root	usr/lib/qubes/qrexec-policy-agent-autostart
lrwxrwxrwx	root/root	usr/lib/qubes/qrexec_client_vm -> ../../bin/qrexec-client-vm
-rwxr-xr-x	root/root	usr/lib/qubes/qubes-rpc-multiplexer
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/system/
-rw-r--r--	root/root	usr/lib/systemd/system/qubes-qrexec-agent.service
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rwxr-xr-x	root/root	usr/share/man/man1/qrexec-client-vm.1.gz
