# qubes-core-agent-linux/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/NetworkManager/
drwxr-xr-x	root/root	etc/NetworkManager/dispatcher.d/
-rwxr-xr-x	root/root	etc/NetworkManager/dispatcher.d/30-qubes-external-ip
-rwxr-xr-x	root/root	etc/NetworkManager/dispatcher.d/qubes-nmhook
drwxr-xr-x	root/root	etc/dhclient.d/
lrwxrwxrwx	root/root	etc/dhclient.d/qubes-setup-dnat-to-ns.sh -> /usr/lib/qubes/qubes-setup-dnat-to-ns
drwxr-xr-x	root/root	etc/init.d/
-rwxr-xr-x	root/root	etc/init.d/qubes-core
-rwxr-xr-x	root/root	etc/init.d/qubes-core-early
-rwxr-xr-x	root/root	etc/init.d/qubes-core-netvm
-rwxr-xr-x	root/root	etc/init.d/qubes-firewall
-rwxr-xr-x	root/root	etc/init.d/qubes-iptables
-rwxr-xr-x	root/root	etc/init.d/qubes-sysinit
-rwxr-xr-x	root/root	etc/init.d/qubes-updates-proxy
-rwxr-xr-x	root/root	etc/init.d/qubes-updates-proxy-forwarder
drwxr-xr-x	root/root	etc/qubes-rpc/
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.Backup
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.ConnectTCP
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.DetachPciDevice
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.Filecopy
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.GetAppmenus
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.GetDate
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.GetImageRGBA
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.InstallUpdatesGUI
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.OpenInVM
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.OpenURL
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.PostInstall
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.RegisterBackupLocation
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.ResizeDisk
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.Restore
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.RestoreById
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.SelectDirectory
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.SelectFile
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.SetDateTime
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.ShowInTerminal
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.StartApp
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.SuspendPost
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.SuspendPostAll
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.SuspendPre
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.SuspendPreAll
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.TemplateDownload
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.TemplateSearch
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.UpdatesProxy
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.VMExec
lrwxrwxrwx	root/root	etc/qubes-rpc/qubes.VMExecGUI -> qubes.VMExec
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.VMRootShell
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.VMShell
-rwxr-xr-x	root/root	etc/qubes-rpc/qubes.WaitForSession
-rw-r--r--	root/root	etc/qubes-suspend-module-blacklist
drwxr-xr-x	root/root	etc/qubes/
drwxr-xr-x	root/root	etc/qubes/applications/
-rw-r--r--	root/root	etc/qubes/applications/README.txt
drwxr-xr-x	root/root	etc/qubes/applications/org.gnome.Terminal.desktop.d/
-rw-r--r--	root/root	etc/qubes/applications/org.gnome.Terminal.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/
-rw-r--r--	root/root	etc/qubes/autostart/README.txt
drwxr-xr-x	root/root	etc/qubes/autostart/deja-dup-monitor.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/deja-dup-monitor.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gcm-apply.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gcm-apply.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gnome-keyring-gpg.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gnome-keyring-gpg.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gnome-keyring-pkcs11.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gnome-keyring-pkcs11.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gnome-keyring-secrets.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gnome-keyring-secrets.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gnome-keyring-ssh.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gnome-keyring-ssh.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gnome-power-manager.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gnome-power-manager.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gnome-screensaver.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gnome-screensaver.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gnome-settings-daemon.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gnome-settings-daemon.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gnome-sound-applet.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gnome-sound-applet.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gpk-update-icon.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gpk-update-icon.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/gsettings-data-convert.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/gsettings-data-convert.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/imsettings-start.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/imsettings-start.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/krb5-auth-dialog.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/krb5-auth-dialog.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/nm-applet.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/nm-applet.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/notify-osd.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/notify-osd.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/orca-autostart.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/orca-autostart.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/org.gnome.SettingsDaemon.XSettings.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/org.gnome.SettingsDaemon.XSettings.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/pulseaudio-kde.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/pulseaudio-kde.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/pulseaudio.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/pulseaudio.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/restorecond.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/restorecond.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/sealertauto.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/sealertauto.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/spice-vdagent.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/spice-vdagent.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/user-dirs-update-gtk.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/user-dirs-update-gtk.desktop.d/30_qubes.conf
drwxr-xr-x	root/root	etc/qubes/autostart/xfce4-notifyd.desktop.d/
-rw-r--r--	root/root	etc/qubes/autostart/xfce4-notifyd.desktop.d/30_qubes.conf
-r--------	root/root	etc/qubes/ip6tables-enabled.rules
-r--------	root/root	etc/qubes/ip6tables.rules
-r--------	root/root	etc/qubes/iptables.rules
drwxr-xr-x	root/root	etc/qubes/post-install.d/
-rwxr-xr-x	root/root	etc/qubes/post-install.d/10-qubes-core-agent-appmenus.sh
-rwxr-xr-x	root/root	etc/qubes/post-install.d/10-qubes-core-agent-features.sh
-rwxr-xr-x	root/root	etc/qubes/post-install.d/90-qubes-core-agent.sh
-rw-r--r--	root/root	etc/qubes/post-install.d/README
drwxr-xr-x	root/root	etc/qubes/rpc-config/
-rw-r--r--	root/root	etc/qubes/rpc-config/qubes.InstallUpdatesGUI
-rw-r--r--	root/root	etc/qubes/rpc-config/qubes.OpenInVM
-rw-r--r--	root/root	etc/qubes/rpc-config/qubes.OpenURL
-rw-r--r--	root/root	etc/qubes/rpc-config/qubes.SelectDirectory
-rw-r--r--	root/root	etc/qubes/rpc-config/qubes.SelectFile
-rw-r--r--	root/root	etc/qubes/rpc-config/qubes.ShowInTerminal
-rw-r--r--	root/root	etc/qubes/rpc-config/qubes.StartApp
-rw-r--r--	root/root	etc/qubes/rpc-config/qubes.VMExecGUI
-rw-r--r--	root/root	etc/qubes/rpc-config/qubes.VMShell+WaitForSession
drwxr-xr-x	root/root	etc/qubes/suspend-post.d/
-rw-r--r--	root/root	etc/qubes/suspend-post.d/README
-rwxr-xr-x	root/root	etc/qubes/suspend-post.d/qvm-sync-clock.sh
drwxr-xr-x	root/root	etc/qubes/suspend-pre.d/
-rw-r--r--	root/root	etc/qubes/suspend-pre.d/README
drwxr-xr-x	root/root	etc/sysconfig/
drwxr-xr-x	root/root	etc/sysconfig/modules/
-rwxr-xr-x	root/root	etc/sysconfig/modules/qubes-core.modules
drwxr-xr-x	root/root	etc/systemd/
drwxr-xr-x	root/root	etc/systemd/system/
-rw-r--r--	root/root	etc/systemd/system/xendriverdomain.service
drwxr-xr-x	root/root	etc/tinyproxy/
-rw-r--r--	root/root	etc/tinyproxy/tinyproxy-updates.conf
-rw-r--r--	root/root	etc/tinyproxy/updates-blacklist
drwxr-xr-x	root/root	etc/xdg/
drwxr-xr-x	root/root	etc/xdg/autostart/
-rw-r--r--	root/root	etc/xdg/autostart/00-qubes-show-hide-nm-applet.desktop
drwxr-xr-x	root/root	etc/xen/
drwxr-xr-x	root/root	etc/xen/scripts/
-rwxr-xr-x	root/root	etc/xen/scripts/vif-qubes-nat.sh
-rwxr-xr-x	root/root	etc/xen/scripts/vif-route-qubes
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/modules-load.d/
-rw-r--r--	root/root	lib/modules-load.d/qubes-core.conf
drwxr-xr-x	root/root	lib/systemd/
drwxr-xr-x	root/root	lib/systemd/system-preset/
-rw-r--r--	root/root	lib/systemd/system-preset/75-qubes-vm.preset
drwxr-xr-x	root/root	lib/systemd/system/
drwxr-xr-x	root/root	lib/systemd/system/ModemManager.service.d/
-rw-r--r--	root/root	lib/systemd/system/ModemManager.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/NetworkManager-wait-online.service.d/
-rw-r--r--	root/root	lib/systemd/system/NetworkManager-wait-online.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/NetworkManager.service.d/
-rw-r--r--	root/root	lib/systemd/system/NetworkManager.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/boot.automount.d/
-rw-r--r--	root/root	lib/systemd/system/boot.automount.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/chronyd.service.d/
-rw-r--r--	root/root	lib/systemd/system/chronyd.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/crond.service.d/
-rw-r--r--	root/root	lib/systemd/system/crond.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/cups-browsed.service.d/
-rw-r--r--	root/root	lib/systemd/system/cups-browsed.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/cups.path.d/
-rw-r--r--	root/root	lib/systemd/system/cups.path.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/cups.service.d/
-rw-r--r--	root/root	lib/systemd/system/cups.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/cups.socket.d/
-rw-r--r--	root/root	lib/systemd/system/cups.socket.d/30_qubes.conf
-rw-r--r--	root/root	lib/systemd/system/dev-xvdc1-swap.service
drwxr-xr-x	root/root	lib/systemd/system/getty@tty.service.d/
-rw-r--r--	root/root	lib/systemd/system/getty@tty.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/org.cups.cupsd.path.d/
-rw-r--r--	root/root	lib/systemd/system/org.cups.cupsd.path.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/org.cups.cupsd.service.d/
-rw-r--r--	root/root	lib/systemd/system/org.cups.cupsd.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/org.cups.cupsd.socket.d/
-rw-r--r--	root/root	lib/systemd/system/org.cups.cupsd.socket.d/30_qubes.conf
-rw-r--r--	root/root	lib/systemd/system/qubes-early-vm-config.service
-rw-r--r--	root/root	lib/systemd/system/qubes-firewall.service
-rw-r--r--	root/root	lib/systemd/system/qubes-iptables.service
-rw-r--r--	root/root	lib/systemd/system/qubes-misc-post.service
-rw-r--r--	root/root	lib/systemd/system/qubes-mount-dirs.service
-rw-r--r--	root/root	lib/systemd/system/qubes-network-uplink.service
-rw-r--r--	root/root	lib/systemd/system/qubes-network-uplink@.service
-rw-r--r--	root/root	lib/systemd/system/qubes-network.service
-rw-r--r--	root/root	lib/systemd/system/qubes-rootfs-resize.service
-rw-r--r--	root/root	lib/systemd/system/qubes-sync-time.service
-rw-r--r--	root/root	lib/systemd/system/qubes-sync-time.timer
-rw-r--r--	root/root	lib/systemd/system/qubes-sysinit.service
-rw-r--r--	root/root	lib/systemd/system/qubes-update-check.service
-rw-r--r--	root/root	lib/systemd/system/qubes-update-check.timer
-rw-r--r--	root/root	lib/systemd/system/qubes-updates-proxy-forwarder.socket
-rw-r--r--	root/root	lib/systemd/system/qubes-updates-proxy-forwarder@.service
-rw-r--r--	root/root	lib/systemd/system/qubes-updates-proxy.service
drwxr-xr-x	root/root	lib/systemd/system/serial-getty@.service.d/
-rw-r--r--	root/root	lib/systemd/system/serial-getty@.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/systemd-random-seed.service.d/
-rw-r--r--	root/root	lib/systemd/system/systemd-random-seed.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/systemd-timesyncd.service.d/
-rw-r--r--	root/root	lib/systemd/system/systemd-timesyncd.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/tinyproxy.service.d/
-rw-r--r--	root/root	lib/systemd/system/tinyproxy.service.d/30_not_needed_in_qubes_by_default.conf
drwxr-xr-x	root/root	lib/systemd/system/tmp.mount.d/
-rw-r--r--	root/root	lib/systemd/system/tmp.mount.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/tor.service.d/
-rw-r--r--	root/root	lib/systemd/system/tor.service.d/30_qubes.conf
drwxr-xr-x	root/root	lib/systemd/system/tor@default.service.d/
-rw-r--r--	root/root	lib/systemd/system/tor@default.service.d/30_qubes.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/qubes-firewall
-rwxr-xr-x	root/root	usr/bin/qubes-open
-rwxr-xr-x	root/root	usr/bin/qubes-vmexec
-rwxr-xr-x	root/root	usr/bin/qvm-copy
lrwxrwxrwx	root/root	usr/bin/qvm-copy-to-vm -> qvm-copy
lrwxrwxrwx	root/root	usr/bin/qvm-move -> qvm-copy
lrwxrwxrwx	root/root	usr/bin/qvm-move-to-vm -> qvm-copy
-rwxr-xr-x	root/root	usr/bin/qvm-open-in-dvm
-rwxr-xr-x	root/root	usr/bin/qvm-open-in-vm
-rwxr-xr-x	root/root	usr/bin/qvm-run-vm
-rwxr-xr-x	root/root	usr/bin/qvm-sync-clock
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/NetworkManager/
drwxr-xr-x	root/root	usr/lib/NetworkManager/conf.d/
-rw-r--r--	root/root	usr/lib/NetworkManager/conf.d/30-qubes.conf
-rw-r--r--	root/root	usr/lib/NetworkManager/conf.d/31-randomize-mac.conf
drwxr-xr-x	root/root	usr/lib/python3.9/
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qubesagent-4.1.26-py3.9.egg-info/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent-4.1.26-py3.9.egg-info/PKG-INFO
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent-4.1.26-py3.9.egg-info/SOURCES.txt
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent-4.1.26-py3.9.egg-info/dependency_links.txt
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent-4.1.26-py3.9.egg-info/top_level.txt
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qubesagent/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__init__.py (EMPTY)
drwxr-xr-x	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/__init__.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/__init__.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/firewall.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/firewall.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/test_firewall.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/test_firewall.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/test_vmexec.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/test_vmexec.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/test_xdg.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/test_xdg.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/vmexec.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/vmexec.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/xdg.cpython-39.opt-1.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/__pycache__/xdg.cpython-39.pyc
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/firewall.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/test_firewall.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/test_vmexec.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/test_xdg.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/vmexec.py
-rw-r--r--	root/root	usr/lib/python3.9/site-packages/qubesagent/xdg.py
drwxr-xr-x	root/root	usr/lib/qubes/
drwxr-xr-x	root/root	usr/lib/qubes/init/
-rwxr-xr-x	root/root	usr/lib/qubes/init/bind-dirs.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/control-printer-icon.sh
-rw-r--r--	root/root	usr/lib/qubes/init/functions
-rwxr-xr-x	root/root	usr/lib/qubes/init/misc-post-stop.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/misc-post.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/mount-dirs.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/network-proxy-setup.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/network-proxy-stop.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/network-uplink-wait.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/qubes-early-vm-config.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/qubes-iptables
-rwxr-xr-x	root/root	usr/lib/qubes/init/qubes-random-seed.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/qubes-sysinit.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/resize-rootfs-if-needed.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/setup-rw.sh
-rwxr-xr-x	root/root	usr/lib/qubes/init/setup-rwdev.sh
-rwxr-xr-x	root/root	usr/lib/qubes/iptables-updates-proxy
-rwxr-xr-x	root/root	usr/lib/qubes/network-manager-prepare-conf-dir
-rwxr-xr-x	root/root	usr/lib/qubes/prepare-suspend
-rwxr-xr-x	root/root	usr/lib/qubes/qfile-agent
-rwsr-xr-x	root/root	usr/lib/qubes/qfile-unpacker
-rwxr-xr-x	root/root	usr/lib/qubes/qopen-in-vm
-rwxr-xr-x	root/root	usr/lib/qubes/qrun-in-vm
-rwxr-xr-x	root/root	usr/lib/qubes/qubes-fix-nm-conf.sh
-rwxr-xr-x	root/root	usr/lib/qubes/qubes-setup-dnat-to-ns
-rwxr-xr-x	root/root	usr/lib/qubes/qubes-sync-clock
-rwxr-xr-x	root/root	usr/lib/qubes/qvm-copy-to-vm.gnome
lrwxrwxrwx	root/root	usr/lib/qubes/qvm-copy-to-vm.kde -> qvm-copy-to-vm.gnome
lrwxrwxrwx	root/root	usr/lib/qubes/qvm-move-to-vm.gnome -> qvm-copy-to-vm.gnome
lrwxrwxrwx	root/root	usr/lib/qubes/qvm-move-to-vm.kde -> qvm-copy-to-vm.gnome
-rwxr-xr-x	root/root	usr/lib/qubes/qvm-template-repo-query
-rwxr-xr-x	root/root	usr/lib/qubes/resize-rootfs
-rwxr-xr-x	root/root	usr/lib/qubes/show-hide-nm-applet.sh
-rwxr-xr-x	root/root	usr/lib/qubes/tar2qfile
-rwxr-xr-x	root/root	usr/lib/qubes/vm-file-editor
-rwxr-xr-x	root/root	usr/lib/qubes/xdg-icon
drwxr-xr-x	root/root	usr/lib/systemd/
drwxr-xr-x	root/root	usr/lib/systemd/user/
drwxr-xr-x	root/root	usr/lib/systemd/user/pulseaudio.service.d/
-rw-r--r--	root/root	usr/lib/systemd/user/pulseaudio.service.d/30_qubes.conf
drwxr-xr-x	root/root	usr/lib/systemd/user/pulseaudio.socket.d/
-rw-r--r--	root/root	usr/lib/systemd/user/pulseaudio.socket.d/30_qubes.conf
drwxr-xr-x	root/root	usr/lib/tmpfiles.d/
-rw-r--r--	root/root	usr/lib/tmpfiles.d/qubes-core-agent-linux.conf
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/qvm-copy.1.gz
-rw-r--r--	root/root	usr/share/man/man1/qvm-open-in-dvm.1.gz
-rw-r--r--	root/root	usr/share/man/man1/qvm-open-in-vm.1.gz
-rw-r--r--	root/root	usr/share/man/man1/qvm-run-vm.1.gz
