# qubes-gui-agent-linux/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/X11/
-rw-r--r--	root/root	etc/X11/Xwrapper.config
drwxr-xr-x	root/root	etc/X11/xinit/
drwxr-xr-x	root/root	etc/X11/xinit/xinitrc.d/
-rwxr-xr-x	root/root	etc/X11/xinit/xinitrc.d/20qt-gnome-desktop-session-id.sh
-rwxr-xr-x	root/root	etc/X11/xinit/xinitrc.d/20qt-x11-no-mitshm.sh
-rwxr-xr-x	root/root	etc/X11/xinit/xinitrc.d/qubes-keymap.sh
-rwxr-xr-x	root/root	etc/X11/xinit/xinitrc.d/z-qubes-session.sh
-rw-r--r--	root/root	etc/X11/xorg-qubes.conf.template
drwxr-xr-x	root/root	etc/init.d/
-rwxr-xr-x	root/root	etc/init.d/qubes-gui-agent
drwxr-xr-x	root/root	etc/pam.d/
-rw-r--r--	root/root	etc/pam.d/qubes-gui-agent
drwxr-xr-x	root/root	etc/profile.d/
-rw-r--r--	root/root	etc/profile.d/qubes-gui.csh
-rw-r--r--	root/root	etc/profile.d/qubes-gui.sh
-rw-r--r--	root/root	etc/profile.d/qubes-session.sh
drwxr-xr-x	root/root	etc/pulse/
-rw-r--r--	root/root	etc/pulse/qubes-default.pa
drwxr-xr-x	root/root	etc/qubes-rpc/
-rw-r--r--	root/root	etc/qubes-rpc/qubes.SetMonitorLayout
drwxr-xr-x	root/root	etc/security/
drwxr-xr-x	root/root	etc/security/limits.d/
-rw-r--r--	root/root	etc/security/limits.d/90-qubes-gui.conf
drwxr-xr-x	root/root	etc/sysconfig/
-rw-r--r--	root/root	etc/sysconfig/desktop
drwxr-xr-x	root/root	etc/sysconfig/modules/
-rwxr-xr-x	root/root	etc/sysconfig/modules/qubes-u2mfn.modules
drwxr-xr-x	root/root	etc/xdg/
-rw-r--r--	root/root	etc/xdg/Trolltech.conf
drwxr-xr-x	root/root	etc/xdg/autostart/
-rw-r--r--	root/root	etc/xdg/autostart/qubes-icon-sender.desktop
-rw-r--r--	root/root	etc/xdg/autostart/qubes-pulseaudio.desktop
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/systemd/
drwxr-xr-x	root/root	lib/systemd/system/
-rw-r--r--	root/root	lib/systemd/system/qubes-gui-agent.service
drwxr-xr-x	root/root	lib/udev/
drwxr-xr-x	root/root	lib/udev/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/70-master-of-seat.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/qubes-change-keyboard-layout
-rwxr-xr-x	root/root	usr/bin/qubes-gui
-rwxr-xr-x	root/root	usr/bin/qubes-gui-runuser
-rwxr-xr-x	root/root	usr/bin/qubes-run-xorg
-rwxr-xr-x	root/root	usr/bin/qubes-session
-rwxr-xr-x	root/root	usr/bin/qubes-set-monitor-layout
-rwxr-xr-x	root/root	usr/bin/start-pulseaudio-with-vchan
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/modules-load.d/
-rw-r--r--	root/root	usr/lib/modules-load.d/qubes-gui.conf
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/xorg-qubes.pc
drwxr-xr-x	root/root	usr/lib/pulse-12.2/
drwxr-xr-x	root/root	usr/lib/pulse-12.2/modules/
-rwxr-xr-x	root/root	usr/lib/pulse-12.2/modules/module-vchan-sink.so
drwxr-xr-x	root/root	usr/lib/qubes/
-rwxr-xr-x	root/root	usr/lib/qubes/icon-sender
drwxr-xr-x	root/root	usr/lib/sysctl.d/
-rw-r--r--	root/root	usr/lib/sysctl.d/30-qubes-gui-agent.conf
drwxr-xr-x	root/root	usr/lib/tmpfiles.d/
-rw-r--r--	root/root	usr/lib/tmpfiles.d/qubes-pulseaudio.conf
-rw-r--r--	root/root	usr/lib/tmpfiles.d/qubes-session.conf
drwxr-xr-x	root/root	usr/lib/xorg/
drwxr-xr-x	root/root	usr/lib/xorg/modules/
drwxr-xr-x	root/root	usr/lib/xorg/modules/drivers/
-rwxr-xr-x	root/root	usr/lib/xorg/modules/drivers/dummyqbs_drv.la
-rwxr-xr-x	root/root	usr/lib/xorg/modules/drivers/dummyqbs_drv.so
-rwxr-xr-x	root/root	usr/lib/xorg/modules/drivers/qubes_drv.so
drwxr-xr-x	root/root	usr/lib/xorg/modules/input/
-rwxr-xr-x	root/root	usr/lib/xorg/modules/input/qubes_drv.la
-rwxr-xr-x	root/root	usr/lib/xorg/modules/input/qubes_drv.so
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/glib-2.0/
drwxr-xr-x	root/root	usr/share/glib-2.0/schemas/
-rw-r--r--	root/root	usr/share/glib-2.0/schemas/20_qubes-gui-vm.gschema.override
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/qubes/
